# PROGRAMS/tech-tip/tech-tip.timer
[Unit]
Description=Run Tech Tip Generator Daily
Requires=tech-tip.service

[Timer]
# Run daily at 8:00 AM
OnCalendar=*-*-* 08:00:00

# Also run on boot if missed
Persistent=true

# Add a small random delay to avoid system load spikes
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
